# ngrok.service.example
[Unit]
Description=ngrok daemon
After=network.target

[Service]
User=<username>
Group=www-data
ExecStart=ngrok http 443
Restart=on-failure

[Install]
WantedBy=multi-user.target